# deploy/transcriber.service
# /etc/systemd/system/transcriber.service
#
# Install:
#   sudo cp deploy/transcriber.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable transcriber
#   sudo systemctl start transcriber
#
# NOTE: Update NODE_BIN path if Node is installed elsewhere (e.g. nvm).
# Find with: which node

[Unit]
Description=Transcriber Next.js App
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/transcriber
Environment=NODE_ENV=production
Environment=PATH=/root/.nvm/versions/node/v24.12.0/bin:/usr/bin:/bin
EnvironmentFile=/var/www/transcriber/.env
ExecStart=/root/.nvm/versions/node/v24.12.0/bin/npm run start
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
